# config/leuitlog-webui.service
[Unit]
Description=LeuitLog Web UI
Documentation=https://github.com/leuitlog/leuitlog
After=network.target leuitlog.service
Wants=leuitlog.service

[Service]
Type=simple
User=leuitlog
Group=leuitlog
WorkingDirectory=/opt/leuitlog
ExecStart=/usr/bin/python3 -m gunicorn --bind 127.0.0.1:8080 --workers 2 --timeout 30 src.leuitlog_webui:app
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=leuitlog-webui

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadOnlyPaths=/var/log/leuitlog

# Resource limits
LimitNOFILE=4096
MemoryMax=128M

[Install]
WantedBy=multi-user.target
